# templates/default/torquebox.conf.erb
description "This is an upstart job file for TorqueBox"

pre-start script
bash << "EOF"
  mkdir -p <%= node[:torquebox][:log_dir] %>
  chown -R torquebox <%= node[:torquebox][:log_dir] %>
EOF
end script

start on started network-services
stop on stopped network-services
respawn

script
bash << "EOF"
  su - torquebox
  <%= @torquebox_dir %>/jboss/bin/standalone.sh <%= @bind_opts %> >> <%= node[:torquebox][:log_dir] %>/torquebox.log 2>&1
EOF
end script